# 打包环境/deb/crond/extract/DEBIAN/postinst
#! /bin/bash

mv /opt/iepCrond/iepCrond.service  /lib/systemd/system/

systemctl daemon-reload && systemctl enable iepCrond.service

