# 打包环境/deb/crond/extract/DEBIAN/prerm
#! /bin/bash

systemctl  stop  iepCrond
systemctl  disable  iepCrond.service
rm -f /lib/systemd/system/iepCrond.service

rm -rf /opt/iepCrond
